# etera-chereta-supervisor.conf
[program:etera-chereta-checker]
command=php /c/xampp/htdocs/etera28082025/artisan etera-chereta:check-expiration
directory=/c/xampp/htdocs/etera28082025
autostart=true
autorestart=true
user=www-data
redirect_stderr=true
stdout_logfile=/c/xampp/htdocs/etera28082025/storage/logs/etera-chereta-checker.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
stopwaitsecs=3600
startsecs=1
startretries=3 